# AI_Brain/infrastructure/systemd/template.service
[Unit]
Description=AI Brain Service Template
After=network.target

[Service]
# Mandatory directives for Service Recovery
Restart=always
RestartSec=10
StartLimitIntervalSec=0

# Application specific configuration
User=ubuntu
WorkingDirectory=/home/ubuntu/ai_brain
ExecStart=/usr/bin/python3 /home/ubuntu/ai_brain/service_script.py
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
